# tailscale/rootfs/etc/s6-overlay/s6-rc.d/magicdns-egress-proxy/run
#!/command/with-contenv bashio
# shellcheck shell=bash
export LOG_FD
# ==============================================================================
# Home Assistant Community App: Tailscale
# Runs the MagicDNS egress proxy
# ==============================================================================

# Note: This service runs when userspace_networking is disabled

# Sets up IPv4 egress DNS proxy for tailscaled/MagicDNS, that forwards queries to IPv4 hassio DNS.

# Note: This script assumes that internal hassio network is available in IPv4.
# In case of an IPv6-only hassio network this IPv4 solution won't work.
# But running an egress DNS proxy on ::1 IPv6 localhost is impossible, because DNS apps can bind to this address also,
# but we can specify only port 53 DNS in resolv.conf for tailscaled.
# So in case of IPv6-only hassio, we can't use [::1]:53 for the egress DNS proxy, but we can't use anything else in resolv.conf for tailscaled.
# For the status of IPv6 support see https://github.com/home-assistant/supervisor/issues/2133

source /usr/lib/trace.sh

readonly DNSMASQ_EGRESS_ADDRESS_IPV4="127.100.100.100"
readonly DNSMASQ_EGRESS_DEFAULT_PORT=53127

# MAGICDNS_MODE is in contenv at /var/run/s6/container_environment
# - "RESTRICTED"   only domains in black-white-list are allowed on ingress side, everything else on egress side
# - "UNRESTRICTED" everything except domains in black-white-list are allowed on ingress side, only domains in black-white-list on egress side
readonly DNSMASQ_BLACK_WHITE_LIST_LOCATION="/run/dnsmasq_black_white_list"

declare hassio_dns_ipv4
declare -a black_list
declare -a white_list
declare domain
declare -a options
declare dnsmasq_egress_port

bashio::log.info "Starting MagicDNS egress proxy..."

function dig_hassio_dns() {
  local type="${1}"
  dig dns.local.hass.io "${type}" +short \
    | { grep -Ev '^;|\.$|^$' || true ;} \
    | head -n 1
}

# Search for unallocated port in case of collision
dnsmasq_egress_port="${DNSMASQ_EGRESS_DEFAULT_PORT}"
while netstat -lntu \
  | awk '{print $4}' \
  | grep -q ":${dnsmasq_egress_port}$"
do
  if ((++dnsmasq_egress_port >= DNSMASQ_EGRESS_DEFAULT_PORT + 1000)); then
    bashio::exit.nok "Can't find free port for MagicDNS egress proxy"
  fi
done

options+=(--no-hosts)
options+=(--no-resolv)
options+=(--conf-file=/dev/null)
options+=(--keep-in-foreground)
options+=(--log-facility='-')
options+=(--cache-size=0)

options+=(--listen-address="${DNSMASQ_EGRESS_ADDRESS_IPV4}")
# on the lo interface this bind-interfaces option has the effect that dnsmasq connects only to the specified address
# and on the lo interface bind-interfaces option doesn't result in a logged warning
# this is intentionally different from the ingress proxy, where on a non-lo interface bind-dynamic option has the same effect
options+=(--bind-interfaces)
options+=(--port="${dnsmasq_egress_port}")

# Hassio DNS's IP addresses
if ! hassio_dns_ipv4=$(dig_hassio_dns A) || \
  bashio::var.is_empty "${hassio_dns_ipv4}"
then
  bashio::exit.nok "Failed to resolve Home Assistant's IPv4 DNS address"
fi

if bashio::var.equals "${MAGICDNS_MODE}" "RESTRICTED"; then
  # This means "Unrestricted" on egress side, ie. everything except domains in black-white-list are allowed

  # Black-list
  readarray -t black_list < "${DNSMASQ_BLACK_WHITE_LIST_LOCATION}"

  # Forward everything to hassio DNS, except the black_list
  bashio::log.info "MagicDNS egress proxy: forwarding everything to ${hassio_dns_ipv4}${black_list:+, except:}"
  options+=(--server=${hassio_dns_ipv4})
  for domain in "${black_list[@]}"; do
    bashio::log.info "  ${domain}"
    options+=(--server=/${domain}/)
  done
elif bashio::var.equals "${MAGICDNS_MODE}" "UNRESTRICTED"; then
  # This means "Restricted" on egress side, ie. only domains in black-white-list are allowed

  # White-list
  readarray -t white_list < "${DNSMASQ_BLACK_WHITE_LIST_LOCATION}"

  # Return NXDOMAIN for everything, except the white_list
  bashio::log.info "MagicDNS egress proxy: returning NXDOMAIN for everything${white_list:+, except forwarding to ${hassio_dns_ipv4} in case of:}"
  options+=(--address=/#/)
  for domain in "${white_list[@]}"; do
    bashio::log.info "  ${domain}"
    options+=(--server=/${domain}/${hassio_dns_ipv4})
  done
else
  bashio::exit.nok "Invalid MAGICDNS_MODE: '${MAGICDNS_MODE}'"
fi

if bashio_custom::trace; then
  options+=(--log-queries)
  options+=(--log-debug)
fi

magicdns-egress-proxy-forwarding setup forwarding "${dnsmasq_egress_port}"

if ! bashio::fs.file_exists "/etc/resolv.for-tailscaled.conf"; then
  # Overwrite it only on the first service startup
  echo "nameserver ${DNSMASQ_EGRESS_ADDRESS_IPV4}" > /etc/resolv.for-tailscaled.conf
fi

# We need to delay the starting of the dependent services until iptables are configured and the conf file is written
echo "" >&3
exec 3>&-

exec dnsmasq "${options[@]}"
