# tailscale/rootfs/etc/s6-overlay/s6-rc.d/magicdns-proxies-reconfigurator/run
#!/command/with-contenv bashio
# shellcheck shell=bash
export LOG_FD
# ==============================================================================
# Home Assistant Community App: Tailscale
# Periodically runs the dnsmasq proxies black/white list initialization
# ==============================================================================

# used by the healthcheck script
# MAGICDNS_PROXIES_RECONFIGURATOR_HEALTH_STATE is in contenv at /var/run/s6/container_environment

# used by the magicdns-ingress-proxy service
readonly MAGICDNS_INGRESS_PROXY_SUPPRESS_FORWARDING_CONFIGURATION_LOCATION="/run/magicdns_ingress_proxy_suppress_forwarding_configuration"
# updated by the magicdns-ingress-proxy service
readonly DNSMASQ_INGRESS_PORT_LOCATION="/run/dnsmasq_ingress_port"

declare current_ipn_dns_config
declare previous_ipn_dns_config=

declare magicdns_proxies_configuration
declare healthy=true

function try {
    set +e
    (set -e; "$@")
    declare -gx TRY_ERROR=$?
    set -e
}

function file_non_empty() {
    local file=${1}

    if [[ -s "${file}" ]]; then
        return "${__BASHIO_EXIT_OK}"
    fi

    return "${__BASHIO_EXIT_NOK}"
}

function restart_proxies() {
  bashio::log.info "MagicDNS proxy configuration has changed"

  bashio::log.info "Disabling incoming DNS queries during reconfiguration..."
  magicdns-ingress-proxy-forwarding setup drop
  magicdns-ingress-proxy-forwarding remove forwarding

  bashio::log.info "Restarting MagicDNS egress proxy..."
  s6-svc -ruwR /run/service/magicdns-egress-proxy

  bashio::log.info "Restarting MagicDNS ingress proxy..."
  touch "${MAGICDNS_INGRESS_PROXY_SUPPRESS_FORWARDING_CONFIGURATION_LOCATION}"
  rm -r -f "${DNSMASQ_INGRESS_PORT_LOCATION}"
  s6-svc -ruwR /run/service/magicdns-ingress-proxy
  rm -r -f "${MAGICDNS_INGRESS_PROXY_SUPPRESS_FORWARDING_CONFIGURATION_LOCATION}"
  if ! file_non_empty "${DNSMASQ_INGRESS_PORT_LOCATION}"; then
    bashio::log.fatal "MagicDNS ingress proxy port is unknown"
    return 1
  fi

  bashio::log.info "Enabling incoming DNS queries after reconfiguration..."
  magicdns-ingress-proxy-forwarding setup forwarding "$(<"${DNSMASQ_INGRESS_PORT_LOCATION}")"
  rm -r -f "${DNSMASQ_INGRESS_PORT_LOCATION}"
  magicdns-ingress-proxy-forwarding remove drop
}

bashio::log.info "Subscribing to IPN message bus..."
# This jq filters the whole state into sg. like:
# {"Resolvers":true,"Routes":["foo.com","tail1234.ts.net.","ts.net."]}
# If the structure changes, that should break the most basic smoke test after a TS client update, so we can fix it before release
/opt/tailscale debug watch-ipn --initial=true --netmap=true --rate-limit=false \
| stdbuf -i0 -oL -eL jq -rcM '{"DNS": .NetMap | has("DNS"), "Resolvers": .NetMap.DNS | has("Resolvers"), "Routes": .NetMap.DNS.Routes // {} | keys}' \
| while IFS= read -r current_ipn_dns_config
do
  if [[ "${current_ipn_dns_config}" == '{"DNS":true,'* && "${current_ipn_dns_config}" != "${previous_ipn_dns_config}" ]]; then

    # This seems to be double check, but not completely:
    # - The ipn watcher will send the initial state, and we compare then the current DNS state
    #   with the previous call of configure-magicdns-proxies from the magicdns-proxies-configurator service.
    #   This way we don't miss any change between the startup of these 2 services.
    # - On subsequent ipn watcher state changes we use the `tailscale dns status` json, that should be stable

    # If we can't check the DNS configuration, the healthcheck script will treat it as an offline situation
    if ! magicdns_proxies_configuration=$(configure-magicdns-proxies test); then
      printf "UNHEALTHY" > /var/run/s6/container_environment/MAGICDNS_PROXIES_RECONFIGURATOR_HEALTH_STATE
      healthy=false
    else
      previous_ipn_dns_config="${current_ipn_dns_config}"
      if bashio::var.false "${healthy}"; then
        printf "HEALTHY" > /var/run/s6/container_environment/MAGICDNS_PROXIES_RECONFIGURATOR_HEALTH_STATE
        healthy=true
      fi
      # Restart dnsmasq proxies if required
      if ! bashio::var.equals "${magicdns_proxies_configuration}" 'unchanged'; then
        try restart_proxies
        if ((TRY_ERROR)); then
          bashio::exit.nok "Failed to reconfigure and restart MagicDNS proxies. Halting app to prevent network loss."
        fi
      fi
    fi
  fi
done
