# tailscale/rootfs/etc/s6-overlay/s6-rc.d/magicdns-proxies-configurator/run
#!/command/with-contenv bashio
# shellcheck shell=bash
export LOG_FD
# ==============================================================================
# Home Assistant Community App: Tailscale
# Runs the dnsmasq proxies black/white list initialization
# ==============================================================================

configure-magicdns-proxies

# Restart the egress proxy, the configuration is changed
bashio::log.info "MagicDNS proxy configuration has changed, restarting MagicDNS egress proxy..."
s6-svc -ruwR /run/service/magicdns-egress-proxy
