# tailscale/rootfs/etc/s6-overlay/s6-rc.d/tailscaled/run
#!/command/with-contenv bashio
# shellcheck shell=bash
# ==============================================================================
# Home Assistant Community Add-on: Tailscale
# Runs tailscale
# ==============================================================================
declare -a options
declare -a routes=()

bashio::log.info 'Starting Tailscale...'

# Remove stale directories/files due to --statedir=/data/state changed in v0.9.0.1
rm -f /data/derpmap.cached.json
rm -rf /data/certs/*
rm -rf /data/files/*
rm -rf /data/ssh/*
[[ -d /data/certs/ ]] && rmdir /data/certs
[[ -d /data/files/ ]] && rmdir /data/files
[[ -d /data/ssh/ ]] && rmdir /data/ssh

options+=(--state=/data/tailscaled.state)
options+=(--statedir=/data/state)

# Opt out of client log upload to log.tailscale.io
options+=(--no-logs-no-support)

# Use userspace networking by default when not set, or when explicitly enabled
if ! bashio::config.has_value "userspace_networking" || \
  bashio::config.true "userspace_networking";
then
  options+=(--tun=userspace-networking)
else
  options+=(--tun=tailscale1)
fi

# In case of non userspace networking add local subnets to ip rules with higher priority than Tailscale's routing
if bashio::config.false "userspace_networking"; then
  readarray -t routes < <(subnet-routes)
  if (( 0 < ${#routes[@]} )); then
      bashio::log.info "Adding local subnets to ip rules with higher priority than Tailscale's routing,"
      bashio::log.info "to prevent routing local subnets if the same subnet is routed within your tailnet."
  fi
  for route in "${routes[@]}"; do
      bashio::log.info "Adding route ${route} to ip rules"
      ip rule add to "${route}" priority 5000 table main
  done
fi

# Run Tailscale
if bashio::debug ; then
  exec /opt/tailscaled "${options[@]}"
else
  bashio::log.notice \
    "Tailscale logs will be suppressed after 200 lines, set add-on's "\
    "configuration option 'log_level' to 'debug' to see further logs"

  /opt/tailscaled "${options[@]}" 2>&1 \
    | stdbuf -i0 -oL -eL \
      sed -n -e '1,200p' \
        -e "201c[further tailscaled logs suppressed, set add-on's configuration option 'log_level' to 'debug' to see further tailscaled logs]"
fi
