# tailscale/rootfs/etc/s6-overlay/s6-rc.d/proxy/run
#!/command/with-contenv bashio
# shellcheck shell=bash
# ==============================================================================
# Home Assistant Community Add-on: Tailscale
# Enables Tailscale Proxy feature
# ==============================================================================

declare domain
declare protocol="http"
declare wait_counter=0

# Check if Tailscale HTTPS is enabled
if ! /opt/tailscale status --self=true --peers=false --json \
  | jq -rce ".CertDomains[0]" > /dev/null;
then
  bashio::log.notice "Tailscale's HTTPS support is disabled, therefore add-on's Tailscale Proxy functionality is disabled"
  bashio::exit.ok
fi

domain=$(/opt/tailscale status --self=true --peers=false --json | jq -rc ".CertDomains[0]")

# Checking if SSL is used
if bashio::var.true "$(bashio::core.ssl)"; then
  bashio::log.warning "Tailscale's HTTPS support is enabled, but Home Assistant is not accessible through plain HTTP connection"
  bashio::log.warning "The add-on will communicate locally with Home Assistant through an encypted HTTPS connection,"
  bashio::log.warning "that unnecessarily uses resources and slows down the communication on a possibly low end device"
  protocol="https+insecure"
fi

# Set up proxy
# Note: Configuration is persistent between reboots, though configuring the same value won't cause duplicate entries,
#       but will overwite previous entry if core.port has been changed
# Note: If a tailnet node is renamed, there will be deprecated entries in status json under .Web.<different-domain:443>,
#       those can be removed only with reseting the serve configuration
if ! /opt/tailscale serve reset \
  || ! /opt/tailscale serve https:443 / "${protocol}://127.0.0.1:$(bashio::core.port)"
then
  bashio::log.error "Unable to configure Tailscale Proxy"
  bashio::exit.nok
fi

# Test proxy configuration
while (( 200 != $(curl -s -o /dev/null -w "%{http_code}" "https://${domain}:443") )); do
  if (( wait_counter++ == 5 )); then
    bashio::log.error "Unable to connect to Home Assistant on ${protocol}://127.0.0.1:$(bashio::core.port)"
    bashio::log.error "Please check your configuration based on the add-on's Documentation under \"Option: proxy\""
    bashio::exit.nok
  fi
  bashio::log.info "Waiting for Tailscale Proxy to be ready..."
  sleep 2
done

bashio::log.info "Tailscale Proxy is enabled:"
bashio::log.info "  Your Home Assistant instance is available within your Tailnet VPN at"
bashio::log.info "  https://${domain}"
