# tailscale/rootfs/etc/s6-overlay/s6-rc.d/init-nginx/run
#!/command/with-contenv bashio
# shellcheck shell=bash
export LOG_FD
# ==============================================================================
# Home Assistant Community Add-on: Tailscale
# Configures NGINX for use with the Tailscale web
# ==============================================================================

echo XXXXXXXXXXXXXXXXXXXXXXXXXXXXXXXXXXXXXXXX $LOG_FD
echo XXXXXXXXXXXXXXXXXXXXXXXXXXXXXXXXXXXXXXXX $__BASHIO_LOG_LEVEL

bashio::var.json \
    interface "$(bashio::addon.ip_address)" \
    port "^$(bashio::addon.ingress_port)"

echo XXXXXXXXXXXXXXXXXXXXXXXXXXXXXXXXXXXXXXXX

bashio::var.json \
    interface "$(bashio::addon.ip_address)" \
    port "^$(bashio::addon.ingress_port)" | cat

echo XXXXXXXXXXXXXXXXXXXXXXXXXXXXXXXXXXXXXXXX

# Generate Ingress configuration
bashio::var.json \
    interface "$(bashio::addon.ip_address)" \
    port "^$(bashio::addon.ingress_port)" \
    | tempio \
        -template /etc/nginx/templates/ingress.gtpl \
        -out /etc/nginx/servers/ingress.conf

cat /etc/nginx/servers/ingress.conf
echo XXXXXXXXXXXXXXXXXXXXXXXXXXXXXXXXXXXXXXXX $?
