# tailscale/rootfs/etc/NetworkManager/dispatcher.d/protect-subnets
#!/command/with-contenv-merge bashio
# shellcheck shell=bash
export LOG_FD
# The shebang 'with-contenv-merge' above is identical with 'with-contenv', but doesn't clear the current environment containing the dispatcher variables

readonly NM_DISPATCHER_LISTENER_QUEUE=/run/nm-dispatcher-listener-queue

# Redirect outputs to the log (&2 is already set up as pipe by s6)
exec 1>&2
if [[ "${LOG_FD:-}" =~ ^[0-9]+$ ]]; then
  eval "exec ${LOG_FD}>&1" || true
fi

case "${NM_DISPATCHER_ACTION}" in
  up|down|dhcp4-change|dhcp6-change)
    bashio::log.debug "Received Network Manager action ${DEVICE_IP_IFACE-} ${NM_DISPATCHER_ACTION}"
    if ! echo "${DEVICE_IP_IFACE-} ${NM_DISPATCHER_ACTION}" > "${NM_DISPATCHER_LISTENER_QUEUE}"; then
      bashio::log.fatal "Failed to send nm-dispatcher action to listener. Halting app to prevent network loss."
      echo -n 1 > /run/s6-linux-init-container-results/exitcode
      exec /run/s6/basedir/bin/halt
    fi
    ;;
  connectivity-change)
    bashio::log.debug "Unhandled Network Manager action ${NM_DISPATCHER_ACTION} ${CONNECTIVITY_STATE-}"
    ;;
  *)
    bashio::log.debug "Unhandled Network Manager action ${DEVICE_IP_IFACE-} ${NM_DISPATCHER_ACTION}"
    ;;
esac
