# tailscale/rootfs/etc/s6-overlay/s6-rc.d/advanced_serve_config/up
/etc/s6-overlay/s6-rc.d/advanced_serve_config/run
